# scripts/supervisord-deploy.conf
[program:deploy]
command=/deploy_api.sh
priority=4
user=root
autorestart=false
autostart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0