# scripts/supervisord-mysql.conf
[program:mysql]
command=mysqld_safe
numprocs=1
autostart=true
autorestart=true
user=root
priority=1