# scripts/supervisord-flyway.conf
[program:flyway]
command=/wait_db_ready.sh
priority=4
user=root
autorestart=false
autostart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0